# supervisor/hailo.conf
[program:hailo]
command=bash -c "ros2 run hailo_rpi_ros2 ros2_node --use-frame --input rpi"
stopsignal=KILL
killasgroup=true
stopasgroup=true
priority=10