# supervisor/hailo.conf
[program:hailo]
command=bash -c "ros2 launch hailo_face_recognition hailo.launch.py"
stopsignal=KILL
killasgroup=true
stopasgroup=true
priority=10