# supervisor/hailo.conf
[program:face_recognition]
command=bash -c "/workspaces/src/hailo_rpi_ros2/supervisor/start_face_recognition.sh"
stopsignal=KILL
killasgroup=true
stopasgroup=true
priority=10